# runtime/vol/kaad/etc/supervisor/conf.d/ad.conf
[program:ad]
# command=service ssh restart
command=/bin/bash /opt/ad/scripts/init.sh