# runtime/vol/kacloud/etc/supervisor/conf.d/nginx.conf
[program:nginx]
# command=service ssh restart
command=/usr/sbin/nginx -g 'daemon off;'
process_name=%(program_name)s
autostart=true
autorestart=true
startretries=5
exitcodes=0,2,70
stopsignal=INT
stopwaitsecs=2
stdout_logfile=/var/log/supervisor/nginx.log