# use-case/mailer/vol/kamail/etc/supervisor/conf.d/sshd.conf
[program:sshd]
# command=service ssh restart
command=/usr/sbin/sshd -D