# runtime/vol/kabind/etc/supervisor/conf.d/named.conf
[program:named]
# command=service ssh restart
command=/usr/sbin/named -g -c /etc/bind/named.conf -u bind